# files/scripts/thinclient/systemd/thinclient-firstboot-admin.service
[Unit]
Description=Create local admin user on first boot
Before=greetd.service
After=systemd-user-sessions.service
Wants=systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/thinclient-firstboot-admin.sh
StandardInput=tty-force
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
