# files/modules/shared-folder/usr/lib/systemd/system/shared-folder.timer
[Unit]
Description=Run shared-folder on boot and hourly

[Timer]
OnBootSec=5min
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
